# deployment.service
[Unit]
After=svnserve.service

[Service]
ExecStart=/usr/local/lib/deploy/execute.sh

[Install]
WantedBy=default.target
